# docs/development/deployment/codex.service.template
[Unit]
Description=Pubkey Quest Codex
After=network.target

[Service]
Type=simple
User=USER
WorkingDirectory=/home/pubkey-test
ExecStart=/home/pubkey-test/codex
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
