# docs/development/deployment/nostr-hero.service.template
[Unit]
Description=Pubkey Quest Game Server
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/pubkey-quest
ExecStart=/path/to/pubkey-quest/pubkey-quest
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables (optional)
# Environment="PORT=8080"

[Install]
WantedBy=multi-user.target
